# deploy/nattvakten-poweroff.path
[Unit]
Description=Watch for Nattvakten power-off requests

[Path]
PathExists=/run/nattvakten/poweroff.request
Unit=nattvakten-poweroff.service

[Install]
WantedBy=multi-user.target
